# tailscale_nodalia_beta/rootfs/etc/s6-overlay/s6-rc.d/support-api/run
#!/command/with-contenv bashio
# shellcheck shell=bash
exec busybox-extras httpd -f -p 127.0.0.1:25910 -u root:root -h /etc/support-api
