# tailscale_nodalia_beta/rootfs/etc/s6-overlay/s6-rc.d/runtime-status/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community Add-on: Tailscale
# Periodically writes runtime status for onboarding diagnostics
# ==============================================================================

readonly RUNTIME_STATE_FILE="/data/tailscale-runtime.json"
readonly POLL_INTERVAL=3

declare status_json
declare backend_state
declare auth_url
declare self_online
declare self_dns_name
declare self_host_name
declare ipv4_list
declare ipv6_list
declare webui_ready
declare now
declare runtime_payload
declare setup_profile
declare webui_readonly
declare share_mode
declare webui_probe_s
declare webui_probe_body
declare webui_probe_payload
declare webui_probe_time
declare webui_http_code
declare direct_webui_url
declare support_json
declare support_enabled
declare support_tailnet_id
declare support_tailnet_match
declare support_eligible
declare support_active
declare support_url
declare support_reason
declare support_expires_at
declare support_ttl_minutes

declare -i ready_streak=0
declare -i start_epoch
declare -i now_epoch
declare -i uptime_sec

start_epoch=$(date +%s)

while true
do
  status_json=$(/opt/tailscale status --json --self=true --peers=false 2>/dev/null || echo '{}')
  backend_state=$(jq -r '.BackendState // "Unknown"' <<< "${status_json}" 2>/dev/null || echo "Unknown")
  auth_url=$(jq -r '.AuthURL // .Self.AuthURL // empty' <<< "${status_json}" 2>/dev/null || true)
  self_online=$(jq -r '.Self.Online // false' <<< "${status_json}" 2>/dev/null || echo "false")
  self_dns_name=$(jq -r '.Self.DNSName // empty' <<< "${status_json}" 2>/dev/null || true)
  self_host_name=$(jq -r '.Self.HostName // empty' <<< "${status_json}" 2>/dev/null || true)

  ipv4_list=$(/opt/tailscale ip -4 2>/dev/null | jq -Rsc 'split("\n") | map(select(length>0))' || echo '[]')
  ipv6_list=$(/opt/tailscale ip -6 2>/dev/null | jq -Rsc 'split("\n") | map(select(length>0))' || echo '[]')

  webui_ready=false
  webui_probe_s="n/a"
  webui_probe_body=""
  webui_http_code="000"
  webui_probe_payload="$(
    curl -sS --max-time 2 -w $'\n__TS_META__:%{time_total}:%{http_code}' "http://127.0.0.1:25899/" 2>/dev/null || true
  )"
  if [[ -n "${webui_probe_payload}" && "${webui_probe_payload}" == *"__TS_META__:"* ]]; then
    webui_probe_time="${webui_probe_payload##*__TS_META__:}"
    webui_http_code="${webui_probe_time##*:}"
    webui_probe_time="${webui_probe_time%:*}"
    webui_probe_body="${webui_probe_payload%$'\n'__TS_META__:*}"
    if [[ -n "${webui_probe_time}" ]]; then
      webui_probe_s="${webui_probe_time}"
    fi
    if [[ "${webui_http_code}" =~ ^(2[0-9][0-9]|3[0-9][0-9]|401|403)$ ]]; then
      if grep -qi "Tailscale web interface is unavailable" <<< "${webui_probe_body}"; then
        webui_ready=false
        ready_streak=0
      else
        webui_ready=true
        ready_streak=$((ready_streak + 1))
      fi
    else
      ready_streak=0
    fi
  else
    ready_streak=0
  fi

  setup_profile=$(bashio::config "setup_profile" "custom")
  if ! bashio::config.has_value "webui_readonly" || bashio::config.true "webui_readonly"; then
    webui_readonly=true
  else
    webui_readonly=false
  fi
  share_mode=$(bashio::config "share_homeassistant" "disabled")
  direct_webui_url=""
  if [[ "${ipv4_list}" != "[]" ]]; then
    direct_webui_url=$(jq -r '.[0] // empty' <<< "${ipv4_list}" 2>/dev/null || true)
    if [[ -n "${direct_webui_url}" ]]; then
      direct_webui_url="http://${direct_webui_url}:5252"
    fi
  fi

  support_json="$(support-tunnel status 2>/dev/null || echo '{}')"
  support_enabled="$(jq -r '.support_enabled // false' <<< "${support_json}" 2>/dev/null || echo "false")"
  support_tailnet_id="$(jq -r '.support_target_id // empty' <<< "${support_json}" 2>/dev/null || true)"
  support_tailnet_match="$(jq -r '.tailnet_match // false' <<< "${support_json}" 2>/dev/null || echo "false")"
  support_eligible="$(jq -r '.eligible // false' <<< "${support_json}" 2>/dev/null || echo "false")"
  support_active="$(jq -r '.active // false' <<< "${support_json}" 2>/dev/null || echo "false")"
  support_url="$(jq -r '.url // empty' <<< "${support_json}" 2>/dev/null || true)"
  support_reason="$(jq -r '.reason // empty' <<< "${support_json}" 2>/dev/null || true)"
  support_expires_at="$(jq -r '.expires_at // empty' <<< "${support_json}" 2>/dev/null || true)"
  support_ttl_minutes="$(jq -r '.ttl_minutes // 0' <<< "${support_json}" 2>/dev/null || echo "0")"

  now_epoch=$(date +%s)
  uptime_sec=$((now_epoch - start_epoch))
  now="$(date -u +"%Y-%m-%dT%H:%M:%SZ")"
  runtime_payload=$(jq -nc \
    --arg updated_at "${now}" \
    --arg backend_state "${backend_state}" \
    --arg auth_url "${auth_url}" \
    --arg self_dns_name "${self_dns_name}" \
    --arg self_host_name "${self_host_name}" \
    --arg setup_profile "${setup_profile}" \
    --arg share_mode "${share_mode}" \
    --arg self_online "${self_online}" \
    --arg webui_ready "${webui_ready}" \
    --arg webui_readonly "${webui_readonly}" \
    --arg webui_probe_s "${webui_probe_s}" \
    --arg direct_webui_url "${direct_webui_url}" \
    --argjson webui_ready_streak "${ready_streak}" \
    --argjson uptime_sec "${uptime_sec}" \
    --arg support_enabled "${support_enabled}" \
    --arg support_tailnet_id "${support_tailnet_id}" \
    --arg support_tailnet_match "${support_tailnet_match}" \
    --arg support_eligible "${support_eligible}" \
    --arg support_active "${support_active}" \
    --arg support_url "${support_url}" \
    --arg support_reason "${support_reason}" \
    --arg support_expires_at "${support_expires_at}" \
    --arg support_ttl_minutes "${support_ttl_minutes}" \
    --arg ipv4 "${ipv4_list}" \
    --arg ipv6 "${ipv6_list}" \
    '{
      updated_at:$updated_at,
      backend_state:$backend_state,
      auth_url:$auth_url,
      self_online: ($self_online == "true"),
      webui_ready: ($webui_ready == "true"),
      webui_readonly: ($webui_readonly == "true"),
      webui_probe_s:$webui_probe_s,
      webui_ready_streak:$webui_ready_streak,
      uptime_sec:$uptime_sec,
      setup_profile:$setup_profile,
      share_mode:$share_mode,
      support_enabled: ($support_enabled == "true"),
      support_tailnet_id:$support_tailnet_id,
      support_tailnet_match: ($support_tailnet_match == "true"),
      support_eligible: ($support_eligible == "true"),
      support_active: ($support_active == "true"),
      support_url:$support_url,
      support_reason:$support_reason,
      support_expires_at:$support_expires_at,
      support_ttl_minutes:($support_ttl_minutes | tonumber? // 0),
      direct_webui_url:$direct_webui_url,
      self_dns_name:$self_dns_name,
      self_host_name:$self_host_name,
      ipv4: ($ipv4 | fromjson? // []),
      ipv6: ($ipv6 | fromjson? // [])
    }')

  printf '%s\n' "${runtime_payload}" > "${RUNTIME_STATE_FILE}"
  sleep "${POLL_INTERVAL}"
done
